# debian/prerm
#!/bin/bash
set -e
if [ -d /run/systemd/system ]; then
    systemctl stop tvheadend.service || true
    systemctl disable tvheadend.service || true
fi
exit 0
